# scripts/framework-relay/com.subbu-ops.framework-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.subbu-ops.framework-relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/subbu_ai_assistant/cortextos/scripts/framework-relay/framework-relay.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/subbu_ai_assistant/.cortextos/default/logs/framework-relay.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/subbu_ai_assistant/.cortextos/default/logs/framework-relay.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/subbu_ai_assistant</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
